# private-vpc/assets/etc/systemd/system/router.service
[Unit]
Description = router setup

[Service]
ExecStart = /usr/local/router/router.sh
Type = oneshot

[Install]
WantedBy = network-online.target
